# deploy/systemd/axiom-publish.timer
[Unit]
Description=Daily after US close - EOD dashboard publish (broker truth)

[Timer]
# 16:30 America/New_York, weekdays only — a weekend run would stamp bogus
# Sat/Sun points onto the equity curve from the broker's stale Friday marks.
OnCalendar=Mon..Fri *-*-* 16:30 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
