# deploy/systemd/axiom-publish.timer
[Unit]
Description=Daily after US close - EOD dashboard publish (broker truth)

[Timer]
# 16:30 America/New_York, every day. Reads the Gateway's account channel
# (NetLiquidation marked at the close) and stamps the canonical daily point.
OnCalendar=*-*-* 16:30 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
