# deploy/systemd/axiom-weights.service
[Unit]
Description=Axiom Tilt - freeze weekly target weights + publish to dashboard
OnFailure=axiom-alert@%n.service

[Service]
Type=oneshot
WorkingDirectory=%h/axiom_tilt_strategy
ExecStart=%h/axiom_tilt_strategy/.venv/bin/python -m trading.run weights
# Publish the frozen target to the dashboard (no Gateway needed):
ExecStartPost=%h/axiom_tilt_strategy/.venv/bin/python -m trading.publish.backfill
