# deploy/systemd/axiom-publish.service
[Unit]
Description=Axiom Tilt - publish dashboard snapshot (audit trail + yfinance, no broker)
OnFailure=axiom-alert@%n.service

[Service]
Type=oneshot
WorkingDirectory=%h/axiom_tilt_strategy
# Broker-free: reconstructs holdings/NAV/equity from trading/audit + yfinance.
ExecStart=%h/axiom_tilt_strategy/.venv/bin/python -m trading.publish
