# deploy/systemd/axiom-publish.timer
[Unit]
Description=Daily after US close - publish dashboard snapshot from the audit trail

[Timer]
# 16:30 America/New_York, every day. The publisher is broker-free (reads the order
# audit + yfinance), so it does not need the IBKR Gateway and has no market-hours guard.
OnCalendar=*-*-* 16:30 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
